# listenbrainz/msb_mapping/docker/crontab
# Create the mapping.
0 4 * * * listenbrainz /usr/local/bin/python /code/mapper/manage.py create-all >> /code/mapper/lb-cron.log 2>&1
