# listenbrainz/mbid_mapping/docker/crontab
# Create the typesense index periodically
0 4 * * * listenbrainz /usr/local/bin/python /code/mapper/manage.py create-all
